# supervisord.conf
[supervisord]
nodaemon=true

[program:pm2]
command:pm2 start /usr/src/app/ecosystem.config.js --env development --no-daemon
startretries:5
